# doc/systemd/monopd.socket
[Unit]
Description=monopd listening socket

[Socket]
ListenStream=1234

[Install]
WantedBy=sockets.target
